# services/dashboard.service
[Unit]
Description=Clawdbot Dashboard
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/clawdbot-dashboard
ExecStart=/usr/bin/node backend/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
